# deploy/emojiopticon.service
[Unit]
Description=Bluesky Emoji Usage Analyzer
After=network.target

[Service]
Type=simple
User=emojiopticon
Group=emojiopticon
EnvironmentFile=/etc/emojiopticon/environment
ExecStart=/opt/emojiopticon/emojiopticon --config /etc/emojiopticon/prompts.toml --port 1337
Restart=always
RestartSec=5

# Security measures
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths=/etc/emojiopticon /var/lib/emojiopticon /var/log/emojiopticon

# Logging
StandardOutput=append:/var/log/emojiopticon/stdout.log
StandardError=append:/var/log/emojiopticon/stderr.log

[Install]
WantedBy=multi-user.target
